# crontab
0 3 1 */2 * /bin/sh /scripts/run_certbot.sh
